# .config/etc/systemd/system/apt-daily.timer.d/override.conf
[Timer]
OnStartupSec=1min
RandomizedDelaySec=1min
