# .config/etc/systemd/system/apt-daily-upgrade.timer.d/override.conf
[Unit]
After=apt-daily.timer

[Timer]
OnStartupSec=2min
RandomizedDelaySec=1min
